# arm工控机配置文件/rc.local
#!/bin/sh -e

/bin/bash /home/jetson/polaris/sip/sip.sh start
/bin/bash /home/jetson/polaris/sip/zlm.sh start

exit 0
~